# packages/pentaho-deploy/pentaho-deploy-0.1.tar.gz/pentaho-deploy-0.1/pentahodeploy/upstart-jobs/pentaho-biserver.conf
# Pentaho BI Server Service

description     "Pentaho BI Server"
author          "Sergio Oliveira <seocam@seocam.com>"

start on started network
stop on stopping network
stop on starting shutdown

exec %(init_script)s 
