# packages/jabbercracky/jabbercracky-0.1.0-py2.6.egg/etc/init.d/jabbercracky
#!/sbin/runscript

opts="depend start stop"

depend() {
	use logger
}

start() {
	ebegin "Starting Jabbercracky daemon"
	start-stop-daemon --start --quiet --pidfile /var/run/jabbercracky.pid --exec /usr/bin/jabbercrackyd.py
	eend $?
}

stop() {
	ebegin "Stopping Jabbercracky daemon"
	start-stop-daemon --stop --quiet --pidfile /var/run/jabbercracky.pid
	eend $?
}
