# packages/baruwa/baruwa-1.1.0.tar.bz2/baruwa-1.1.0/extras/packaging/baruwa.cron.monthly
#!/bin/sh
#
# Baruwa
#
# Send PDF reports
baruwa-admin sendpdfreports