# packages/silk-deployment/silk-deployment-0.1.1.tar.gz/silk-deployment-0.1.1/silk/cfg_templates/supervisord.conf
[program:%(site)s]
command=%(cmd)s
directory=%(remote_root)s
user=nobody
autostart=true
autorestart=true
stdout_logfile=%(remote_root)s/logs/supervisor.log
redirect_stderr=True
environment=%(process_env)s
